# playbooks/roles/gammu-smsd/templates/watcher.timer.j2
[Unit]
Description=Watcher Timer

[Timer]
OnCalendar={{ gammu_watcher_cron }}

[Install]
WantedBy=timers.target
